# cronjob-dehydrated.service
[Unit]
Description=letsencrypt/acme client cronjob

[Service]
ExecStart=/usr/sbin/dehydrated -c
StandardOutput=null
