# cronjob-dehydrated.timer
[Unit]
Description=letsencrypt/acme client cronjob timer

[Timer]
OnCalendar=Tue *-*-* 02:42:00
Persistent=true

[Install]
WantedBy=cronjobs.target
